# Vagrant-Ansible/ANSIBLE_NOVAGRANT_PROVISION/roles/backup/files/backup_copy.service
[Unit]
Description=Create Backup with copy.sh

[Service]
Type=simple
ExecStart=/home/johnd/copy.sh
User=johnd
Group=johnd
Restart=always

[Install]
WantedBy=multi-user.target
